# dot_config/systemd/user/auto-lock.service.tmpl
[Unit]
Description=Auto lock X session
BindsTo=i3-session.target

[Service]
Type=simple
Restart=always
ExecStart={{ .chezmoi.homeDir }}/.local/bin/auto-lock

[Install]
WantedBy=i3-session.target
