# dot_config/systemd/user/dresrv-online.target
[Unit]
Description=DreSRV is available
Wants=dresrv-gpg.service
Wants=dresrv-sshfs.service
Wants=dresrv-tunnel.service
