# dot_config/systemd/user/nextcloud.service
[Unit]
Description=Nextcloud client
Documentation=man:nextcloud(1)
BindsTo=graphical-session.target

[Service]
Type=forking
ExecStart=/usr/bin/nextcloud --background

[Install]
WantedBy=graphical-session.target
