# dot_config/systemd/user/swayidle.service.tmpl
[Unit]
Description=Swayidle daemon
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart={{ .chezmoi.homeDir }}/.local/bin/swayidle

[Install]
WantedBy=sway-session.target
