# dot_config/systemd/user/ws-online.target
[Unit]
Description=Workstation is available
BindsTo=default.target
Wants=ws-gpg.service
Wants=ws-proxy.service
