# dot_config/systemd/user/ws-online.target
[Unit]
Description=Workstation is available
